# tests/render-supervisor-group/fixtures/test-app.conf

[program:web]
command=env web.sh
autostart=true
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog
environment=KEY2='VALUE2', KEY1='VALUE1', PORT=5000
directory=/opt/test-app/
user=test-app

[program:worker-1]
command=env worker.sh
autostart=true
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog
environment=KEY2='VALUE2', KEY1='VALUE1'
directory=/opt/test-app/
user=test-app

[program:worker-2]
command=env worker.sh
autostart=true
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog
environment=KEY2='VALUE2', KEY1='VALUE1'
directory=/opt/test-app/
user=test-app

[group:test-app]
programs=web, worker-1, worker-2
