# templates/supervisor-group.conf.erb
<% node['procs'].each do |p| %>
[program:<%= p[:name] %>]
command=env <%= p[:command] %>
autostart=true
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog
environment=<%= node['env'].sort.map{|k,v| "#{k}=\"#{v}\""}.join(',') %>
directory=<%= @cfg.cwd %>
user=<%= @cfg.user %>
<% end %>

[group:<%= @cfg.app_name %>]
programs=<%= node['procs'].map{|e| e['name'] }.join(",") %>
